# archlive/airootfs/etc/systemd/system/tor-watchdog.timer
[Unit]
Description=ArchNesic — periodic Tor health check
After=tor.service
Requires=tor.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tor-watchdog.sh

[Timer]
OnBootSec=30
OnUnitActiveSec=30
Unit=tor-watchdog.service

[Install]
WantedBy=timers.target
