# archlive/airootfs/etc/systemd/system/tor.service.d/no-persist.conf
[Service]
ExecStartPre=-/usr/bin/rm -rf /var/lib/tor
ExecStartPre=/usr/bin/mkdir -p /var/lib/tor
ExecStartPre=/usr/bin/mount -t tmpfs -o size=16M,noexec,nosuid,nodev tmpfs /var/lib/tor
